# scripts/dev-watch.service.example
# systemd user service for CrypRQ Development Watcher
# Copy to ~/.config/systemd/user/dev-watch.service
# Replace /path/to/your/repo with your actual repository path

[Unit]
Description=CrypRQ Live Green Dev Watcher
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/path/to/your/repo
EnvironmentFile=%h/path/to/your/repo/.env.devwatch
ExecStart=/usr/bin/env bash scripts/dev-watch.sh
Restart=on-failure
RestartSec=3
# Avoid multiple instances (your script already locks, this is extra)
RuntimeMaxTasks=10000

[Install]
WantedBy=default.target

